# Custom/services/voice-control.service
[Unit]
Description=Aranea Voice Control
After=aranea-server.service
Requires=aranea-server.service

[Service]
WorkingDirectory=/home/piotr/Aranea_code
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /home/piotr/Aranea_code/voice_control.py
Restart=always
RestartSec=2
User=piotr
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target